# conf.d/websockify.conf
[program:websockify]
command=websockify --web /usr/share/novnc %(ENV_NOVNC_PORT)s %(ENV_VNC_SERVER)s
environment=NOVNC_PORT=%(ENV_NOVNC_PORT)s,VNC_SERVER=%(ENV_VNC_SERVER)s
autorestart=true